# Airflow_dbt/airflow_ubun_psql/airflow.ini
[program:airflow-webserver]
command=su - airflow -c "airflow webserver"
directory=/usr/local/airflow
autostart=true
autorestart=true
startretries=3
stderr_logfile=/usr/local/airflow/logs/webserver.err.log
stdout_logfile=/usr/local/airflow/logs/webserver.log
user=root

[program:airflow-sch]
command=su - airflow -c "airflow scheduler" 
directory=/usr/local/airflow
autostart=true
autorestart=true
startretries=3
stderr_logfile=/usr/local/airflow/logs/scheduler.err.log
stdout_logfile=/usr/local/airflow/logs/scheduler.log
user=root

[program:airflow-worker]
command=su - airflow -c "airflow celery worker" 
directory=/usr/local/airflow
autostart=true
autorestart=true
startretries=3
stderr_logfile=/usr/local/airflow/logs/worker.err.log
stdout_logfile=/usr/local/airflow/logs/worker.log
user=root
